# templates/netcontrol.cron

PATH="$PATH:/usr/sbin:/usr/local/bin/:/usr/bin:/etc/netControl/scripts/:/etc/netControl/scripts/include/"
* * * * * root /etc/netControl/scripts/minutely.sh
